# bird6-elsa/files/bird6-elsa.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2010-2012 OpenWrt.org

START=50

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	( SERVICE_MATCH_NAME=1 service_start /usr/sbin/bird6-elsaloop )
#	service_start /usr/sbin/bird6-elsa -d
}

stop() {
	( SERVICE_MATCH_NAME=1 service_stop /usr/sbin/bird6-elsaloop )
#	service_stop /usr/sbin/bird6-elsa
}

reload() {
	service_reload /usr/sbin/bird6-elsa
}
